# templates/default/systemd.service.erb
[Unit]
Description = GlassFish instance
After = syslog.target network.target remote-fs.target mysqld.service
StartLimitIntervalSec=0

[Install]
WantedBy = multi-user.target

[Service]
User=<%= node['glassfish']['user'] %>
ExecStart = <%= "#{@start_domain_command}" %>
ExecStop = <%= "#{@stop_domain_command}" %>
ExecReload = <%= "#{@restart_domain_command}" %>
Type = forking
TimeoutStartSec = <%= @start_domain_timeout %>
TimeoutStopSec = <%= @stop_domain_timeout %>
Restart=always
RestartSec=5