# templates/default/systemd.service.erb
<% prefix = @authbind ? '/usr/bin/authbind --deep ' : '' %>
[Unit]
Description = GlassFish instance
After = syslog.target network.target remote-fs.target mysqld.service <%= node['install']['systemd']['after'] %>
StartLimitIntervalSec=0

[Install]
WantedBy = multi-user.target

[Service]
User=<%= node['glassfish']['user'] %>
ExecStart = <%= "#{prefix} #{@start_domain_command}" %>
ExecStop = <%= "#{prefix} #{@stop_domain_command}" %>
ExecReload = <%= "#{prefix} #{@restart_domain_command}" %>
Type = forking
TimeoutStartSec = <%= @start_domain_timeout %>
TimeoutStopSec = <%= @stop_domain_timeout %>
Restart=always
RestartSec=5